# files/system/etc/skel/.config/autostart/vicinae.desktop
[Desktop Entry]
Type=Application
Name=Vicinae
Exec=vicinae server
X-GNOME-Autostart-enabled=true
NoDisplay=true
