# files/system/etc/xdg/autostart/buzz.zephyos.hello.desktop
[Desktop Entry]
Type=Application
Name=ZephyrOS Hello
Exec=flatpak run buzz.zephyos.hello
X-GNOME-Autostart-enabled=true
NoDisplay=true
